# scripts/systemd/psid.service.in
#
# ParaStation
#
# Copyright (C) 2016, 2017 ParTec Cluster Competence Center GmbH, Munich
# Copyright (C) 2022 ParTec AG, Munich
#
# This file may be distributed under the terms of the Q Public License
# as defined in the file LICENSE.QPL included in the packaging of this
# file.
#
[Unit]
Description=The ParaStation management Daemon
Wants=local-fs.target network-online.target
After=local-fs.target network-online.target

[Service]
Type=forking
IgnoreSIGPIPE=no
ExecStart=@bindir@/psid
Delegate=yes

[Install]
WantedBy=multi-user.target
