# scripts/xinetd.d/psidstarter.in
# default: on
# description: ParaStation management daemon.
service psid
{
	socket_type	= stream
	protocol	= tcp
	type		= UNLISTED
	port		= 888
	wait		= no
	user		= root
	server		= @sbindir@/psid
	log_on_failure	+= USERID
	disable		= no
}
